# Configurazioni/Controllo missione (Edificio C)/2°DNS/bind/rc.local

ifconfig eht0 192.168.4.2 netmask 255.255.255.0 broadcast 192.168.4.255

route add default gw 192.168.4.254